# docker/rootfs/etc/s6-overlay/s6-rc.d/init-search-index/run
#!/command/with-contenv bashio
# shellcheck shell=bash

local -r index_version=3
local -r index_version_file=${DATA_DIR}/.index_version

if [[ (! -f "${index_version_file}") || $(<"${index_version_file}") != "$index_version" ]]; then
	bashio::log.info "Search index out of date. Updating..."
	python3 manage.py document_index reindex --no-progress-bar
	echo ${index_version} | tee "${index_version_file}" >/dev/null
fi
