# docker/rootfs/etc/s6-overlay/s6-rc.d/init-start/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "paperless-ngx docker container starting..."
