# docker/rootfs/etc/s6-overlay/s6-rc.d/init-superuser/run
#!/command/with-contenv bashio
# shellcheck shell=bash

if [[ -n "${PAPERLESS_ADMIN_USER}" ]]; then
	python3 manage.py manage_superuser
fi
