# docker/rootfs/etc/s6-overlay/s6-rc.d/init-superuser/up
s6-setuidgid paperless /etc/s6-overlay/s6-rc.d/init-superuser/run
