# docker/rootfs/etc/s6-overlay/s6-rc.d/init-modify-user/run
#!/command/with-contenv bashio
# shellcheck shell=bash

set -eu

local -r usermap_original_uid=$(id -u paperless)
local -r usermap_original_gid=$(id -g paperless)
local -r usermap_new_uid=${USERMAP_UID:-$usermap_original_uid}
local -r usermap_new_gid=${USERMAP_GID:-${usermap_original_gid:-$usermap_new_uid}}
if [[ ${usermap_new_uid} != "${usermap_original_uid}" || ${usermap_new_gid} != "${usermap_original_gid}" ]]; then
	bashio::log.info "Mapping UID and GID for paperless:paperless to $usermap_new_uid:$usermap_new_gid"
	usermod -o -u "${usermap_new_uid}" paperless
	groupmod -o -g "${usermap_new_gid}" paperless
fi

bashio::log.info "───────────────────────────────────────"
bashio::log.info "GID/UID"
bashio::log.info "───────────────────────────────────────"
bashio::log.info "
User UID:    $(id -u paperless)
User GID:    $(id -g paperless)
───────────────────────────────────────
"
