# docker/rootfs/etc/s6-overlay/s6-rc.d/init-system-checks/run
#!/command/with-contenv bashio
# shellcheck shell=bash

# Explicitly run the Django system checks
bashio::log.info "Running Django checks"
python3 manage.py check
